# files/supervisor.conf
[supervisord]
nodaemon=true

[program:app]
priority=10
directory=/usr/local/crashplan
command=/opt/CrashPlan.sh
user=nobody
autorestart=true
autostart=true
